# hiperium-identity/nodes-config/wildfly-upstart.conf
description "Start the Wildfly Server on boot time in Ubuntu distributions."
author "Andres Solorzano"

start on runlevel [2345]
stop on shutdown

exec service wildfly start
